# templates/vaultwarden.ini
[program:vaultwarden]
directory = %(ENV_HOME)s/vaultwarden/output/
command = %(ENV_HOME)s/vaultwarden/output/vaultwarden
autostart = yes
autorestart = yes
startsecs = 60
environment = ENV_FILE="%(ENV_HOME)s/vaultwarden/.env",DATA_FOLDER="%(ENV_HOME)s/vaultwarden/data"
